# glusterfs_client_confd.conf
[program:glusterfs_client_confd]
command=/bin/bash -c "confd -node ${COREOS_PRIVATE_IPV4}:4001 -config-file /glusterfs/confd.toml"
redirect_stderr=true
priority=999
